# etc/littlebrother-supervisord.ini
[program:littlebrother]
command = /usr/local/bin/uwsgi --ini /srv/etc/uwsgi/littlebrother.ini
user = django
